# etc/rc.local
#!/bin/bash
bash /etc/rc.local.real

exit 0
